# run-your-own-dojo/base-infrastructure/scripts/request.service
[Unit]
Description=Request Password Service

[Service]
Type=simple
User=ec2-user
ExecStart=/usr/bin/request
Restart=always

[Install]
WantedBy=multi-user.target
